# systemd/cad_listener.service
[Unit]
Description=RunSheet CAD Listener - glenmoorefc
After=network.target postgresql.service runsheet.service
Wants=runsheet.service

[Service]
Type=simple
User=dashboard
WorkingDirectory=/opt/runsheet/cad
ExecStart=/opt/runsheet/runsheet_env/bin/python cad_listener.py \
    --port 19117 \
    --tenant glenmoorefc \
    --api-url https://glenmoorefc.cadreport.com \
    --timezone America/New_York
Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cad_listener_glenmoorefc

[Install]
WantedBy=multi-user.target
